# system_files/usr/lib/systemd/user/xdg-desktop-portal-gnome.service.d/10-wait-for-niri.conf
# Close the screencast-portal startup race under niri.
#
# niri serves screencasting through GNOME's Mutter ScreenCast D-Bus interface,
# which xdg-desktop-portal-gnome probes exactly once at startup (with
# DO_NOT_AUTO_START) and never retries. On a session restart -- logout/login, not
# a full reboot -- the portal can D-Bus-reactivate ahead of the new niri, fail
# that probe, and latch into "Non-compatible display server, exposing settings
# only" mode: it silently drops org.freedesktop.impl.portal.ScreenCast and screen
# sharing fails until the portal is restarted by hand. Blocking here until niri
# owns org.gnome.Mutter.ServiceChannel before the backend initializes closes the
# race.
#
# Guarded on XDG_CURRENT_DESKTOP so it is a no-op outside niri: KDE routes every
# portal to the kde backend (kde-portals.conf) and never activates this service,
# but the guard makes that explicit and avoids ever blocking on a name KWin does
# not provide. ExecStart always runs; on the normal path niri already owns the
# name and the loop returns on its first iteration, so this costs nothing.
#
# The dollar sign is doubled ($$) because systemd performs its own $VAR expansion
# on Exec lines; $$ passes a literal $ through to bash, which does the real work.
[Service]
ExecStartPre=/usr/bin/bash -c 'case "$$XDG_CURRENT_DESKTOP" in *niri*) for ((i=0;i<100;i++)); do busctl --user status org.gnome.Mutter.ServiceChannel >/dev/null 2>&1 && exit 0; sleep 0.1; done ;; esac; exit 0'
