# system_files/usr/lib/systemd/user/niri.service.d/10-reap-secrets.conf
# Keep KWallet working across a manual niri restart.
#
# KWallet's Secret Service (org.freedesktop.secrets) is served by ksecretd, which
# plasmalogin's PAM stack starts per login inside the logind session scope. With
# logind's KillUserProcesses=no, that ksecretd is not reaped when the session
# ends; it lingers in a "closing" session and makes the next login's ksecretd lose
# the bus-name race, after which a password-less, D-Bus-activated ksecretd hijacks
# the name and cannot unlock the wallet. tishy-reap-session-secrets reaps the
# ending session's ksecretd (and any password-less hijacker) so the next login
# starts on a clean bus. See that script for the full rationale.
#
# ExecStopPost runs after niri exits, while XDG_SESSION_ID still names the ending
# session. niri.service is Restart=no, so a fresh plasmalogin login always follows
# and respawns ksecretd via PAM -- this never strands a session without a Secret
# Service. The leading '-' keeps a transient non-zero exit (e.g. a process that
# vanished between scan and kill) from being recorded as a unit failure.
[Service]
ExecStopPost=-/usr/local/libexec/tishy-reap-session-secrets
